# systemd/provision-ble.service
[Unit]
Description=BLE Wi-Fi Provisioning Service
After=bluetooth.service NetworkManager.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/provision-ble
Restart=on-failure
RestartSec=2
User=root
Group=root

# Ensure BLE + D-Bus access
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
